# chime.service
[Unit]
Description=Campus Evening Chime System
After=network.target sound.target

[Service]
Type=simple
User=pi
Group=pi
# ワーキングディレクトリは git clone した先に合わせる
WorkingDirectory=/home/pi/steam5pm
# Pythonのパスと実行ファイルを指定 (main.py -> run.py)
ExecStart=/usr/bin/python3 /home/pi/steam5pm/run.py
# 自動復旧設定
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
